# etc/udev/rules.d/90-bspwm-hotplug-monitor.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/denied/.Xauthority", RUN+="/usr/local/bin/bspwm-xrandr.sh"
